# .vps/systemd/copypatrol-api-gunicorn.socket
[Unit]
Description=copypatrol api gunicorn socket

[Socket]
ListenStream=/run/copypatrol-api-gunicorn.sock

[Install]
WantedBy=sockets.target
